# deploy/varagh.service
[Unit]
Description=Varagh game server
After=network.target

[Service]
Type=simple
User=varagh
WorkingDirectory=/opt/varagh
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=HOST=127.0.0.1
Environment=WEB_ORIGIN=https://varagh.artanova.top
Environment=DATABASE_URL=/opt/varagh/data/varagh.db
Environment=NODE_OPTIONS=--disable-warning=ExperimentalWarning
ExecStart=/usr/local/bin/node /opt/varagh/server.cjs
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
